# tests/client/supervisord.conf

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:ipsec]
command=/usr/sbin/starter --nofork
autostart=true
autorestart=true
stderr_logfile=/var/log/ipsec.err.log
stdout_logfile=/var/log/ipsec.out.log

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/var/log/sshd.err.log
stdout_logfile=/var/log/sshd.out.log

[program:zabbix_agent]
command=/usr/sbin/zabbix_agentd -c /etc/zabbix/zabbix_agentd.conf
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/zabbix-agent.err.log
stdout_logfile=/var/log/zabbix-agent.out.log

[program:iperf3]
command=/usr/bin/iperf3 -s -D
autostart=true
autorestart=true
stderr_logfile=/var/log/iperf3.err.log
stdout_logfile=/var/log/iperf3.out.log
startretries=3
