# playbooks/files/strongswan.service
[Unit]
Description=IPsec IKEv1/IKEv2 daemon
After=network.target

[Service]
ExecStart=/usr/libexec/ipsec/starter --nofork
ExecStop=/usr/libexec/ipsec/starter stop
ExecReload=/usr/libexec/ipsec/starter reload
Type=forking
Restart=on-failure

[Install]
WantedBy=multi-user.target